# templates/default/machine/init.d.erb
#!/bin/sh
#
# <%= @init_script_name %>
# Subsystem file for AppDynamics machine agent
#
# chkconfig: 2345 95 05 (1)
# description: AppDynamics machine agent

# source function library
if [ -r /lib/lsb/init-functions ]; then
    . /lib/lsb/init-functions
elif [ -r /etc/rc.d/init.d/functions ]; then
    . /etc/rc.d/init.d/functions
elif [ -r /etc/init.d/functions ]; then
    . /etc/init.d/functions
fi

INSTALL_DIR="<%= @install_dir %>"
RETVAL=0
PROG="$INSTALL_DIR/run.sh"
pid_file="<%= @pid_file %>"

start() {
    echo -n $"Starting AppDynamics machine agent:"
    daemon --pidfile=${pid_file} ${PROG}
    RETVAL=$?
    echo
}

stop() {
    echo -n $"Stopping AppDynamics machine agent:"
    killproc -p ${pid_file} "$PROG" -TERM
    RETVAL=$?
    [ $? -eq 0 ] && rm -f ${pid_file}
    echo
}

is_running() {
    [ -f "$pid_file" ] && ps -p `get_pid` > /dev/null 2>&1
}
get_pid() {
    cat "$pid_file"
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    status)
       if is_running(); then 
         echo "Running"
       else
         echo "Not running"
         exit 1
       fi
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        RETVAL=1
        ;;
esac

exit $RETVAL
